# example/bt-audio-gateway/dist/bt-audio-gateway-0.2.0/systemd/bt-audio-update.service
[Unit]
Description=Update bt-audio-gateway from release manifest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/bt-audio-gateway/bt-audio-gateway.env
ExecStart=/opt/bt-audio-gateway/current/bin/update.sh
