# example/dvd-ingester/dist/dvd-ingester-1.0.2/systemd/dvd-ingester-doctor.timer
[Unit]
Description=Run DVD ingest doctor periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=dvd-ingester-doctor.service
Persistent=true

[Install]
WantedBy=timers.target
