# example/dvd-ingester/dist/dvd-ingester-0.1.9/systemd/dvd-ingester-doctor.service
[Unit]
Description=Run dvd-ingester health checks

[Service]
Type=oneshot
EnvironmentFile=/etc/dvd-ingester/dvd-ingester.env
ExecStart=/opt/dvd-ingester/current/bin/doctor.sh
