# example/dvd-ingester/systemd/dvd-ingester-doctor.timer
[Unit]
Description=Periodic dvd-ingester health check

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
