# example/bt-audio-gateway/dist/bt-audio-gateway-0.2.2/systemd/bt-audio-watch.service
[Unit]
Description=Bluetooth audio device autoconnect watcher
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
EnvironmentFile=/etc/bt-audio-gateway/bt-audio-gateway.env
ExecStart=/opt/bt-audio-gateway/current/bin/bt-audio-watch
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
