# example/dvd-ingester/dist/dvd-ingester-0.3.0/systemd/dvd-rip@.service
[Unit]
Description=Auto-rip DVD from /dev/%I
After=network-online.target remote-fs.target
Wants=network-online.target
RequiresMountsFor=/mnt/nas

[Service]
Type=oneshot
EnvironmentFile=/etc/dvd-ingester/dvd-ingester.env
ExecStart=/opt/dvd-ingester/current/bin/dvd-rip-one /dev/%I
TimeoutStartSec=12h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
