# example/dvd-ingester/dist/dvd-ingester-1.0.3/systemd/dvd-ingester-doctor.service
[Unit]
Description=DVD ingest doctor check
Documentation=file:/opt/dvd-ingester/current/doc/README.md

[Service]
Type=oneshot
EnvironmentFile=-/etc/dvd-ingester/dvd-ingester.env
ExecStart=/opt/dvd-ingester/current/bin/doctor.sh
PrivateTmp=true
NoNewPrivileges=true
