# example/dvd-ingester/dist/dvd-ingester-1.0.3/systemd/dvd-ingester-update.timer
[Unit]
Description=Run DVD ingest update check periodically

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
Unit=dvd-ingester-update.service
Persistent=true

[Install]
WantedBy=timers.target
