# example/dvd-ingester/dist/dvd-ingester-1.0.1/systemd/dvd-ingester-ha-mqtt.timer
[Unit]
Description=Refresh DVD ingest Home Assistant MQTT state

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=dvd-ingester-ha-mqtt.service
Persistent=true

[Install]
WantedBy=timers.target
