# example/dvd-ingester/dist/dvd-ingester-1.0.1/systemd/dvd-ingester-update.service
[Unit]
Description=DVD ingest release update check
Documentation=file:/opt/dvd-ingester/current/doc/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/dvd-ingester/dvd-ingester.env
ExecStart=/opt/dvd-ingester/current/bin/update.sh
PrivateTmp=true
NoNewPrivileges=true
